# set-volume.service
[Unit]
Description=Set Audio Volume to Maximum
After=sound.target
Wants=sound.target

[Service]
Type=oneshot
User=pi
ExecStart=/bin/bash -c '/usr/bin/amixer -c 0 sset Playback 100%% || true'
ExecStart=/bin/bash -c '/usr/bin/amixer -c 0 sset Headphone 100%% || true'
ExecStart=/bin/bash -c '/usr/bin/amixer -c 0 sset Speaker 100%% || true'
ExecStart=/bin/bash -c '/usr/bin/amixer -c 0 sset Capture 100%% || true'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
